# supervisord.conf
[program:postgres]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/14/main
user=postgres
priority=10
autorestart=true
stopsignal=INT
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.err.log

[program:redis]
command=/opt/redis-stack/bin/redis-stack-server
    --dir /var/lib/redis-stack
    --appendonly yes
    --protected-mode no
    --maxmemory-policy allkeys-lru
user=redis
priority=20
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.err.log

[program:nmocore]
command=/usr/local/bin/uvicorn app:app
    --host 127.0.0.1 --port 8002
    --h11-max-incomplete-event-size 16777216
    --log-level info
directory=/app
priority=30
autorestart=true
startretries=50
startsecs=3
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/var/log/supervisor/nmocore.log
stderr_logfile=/var/log/supervisor/nmocore.err.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=40
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err.log
